# deploy/ec2/salary-management.service
[Unit]
Description=Salary Management API stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/salary-management
Environment=DEPLOY_ENV_DIR=/etc/salary-management
ExecStart=/usr/bin/docker compose -f deploy/ec2/docker-compose.ec2.yml up -d --build
ExecStop=/usr/bin/docker compose -f deploy/ec2/docker-compose.ec2.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
